# meta-tgp-core/recipes-varserver/varserver/files/varserver.service
[Unit]
Description=varserver service

[Service]
ExecStart=/usr/bin/varserver

[Install]
WantedBy=multi-user.target
Alias=varserver.service

